# systemd/aw5d-lcd.service
[Unit]
Description=iBUYPOWER AW5D cooler LCD (CPU gauge from lm-sensors)
Documentation=https://github.com/claygorman/aw5d-linux
After=graphical-session.target

[Service]
Type=simple
# %h expands to the user's home. install.sh drops the script here.
# Optional config: set AW5D_INTERVAL (seconds between frames) in ~/.config/aw5d-lcd.env.
# The '-' makes the file optional; default interval is 1s if it's absent.
EnvironmentFile=-%h/.config/aw5d-lcd.env
ExecStart=/usr/bin/python3 %h/.local/share/aw5d-lcd/aw5d_lcd.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
